# crontabs
# m h  dom mon dow   command
MAILTO=lucas@debian.org
UAR=/org/udd.debian.org/update-and-run.sh
# Debian Sources/Packages
0 */6 * * * $UAR debian-squeeze debian-sid debian-experimental
0 3 * * * $UAR debian-etch debian-backports-etch debian-volatile-etch debian-etch-proposed-updates debian-etch-security debian-lenny debian-backports-lenny debian-volatile-lenny debian-lenny-proposed-updates debian-lenny-security debian-squeeze-proposed-updates debian-squeeze-security
# Ubuntu Sources/Packages
0 1 * * * $UAR ubuntu-jaunty
15 2 */2 * * $UAR ubuntu-hardy ubuntu-intrepid
# Various simple things
0 2 * * * $UAR debian-popcon ubuntu-popcon upload-history orphaned_packages lintian debtags
15 2 * * * $UAR carnivore
0 */12 * * * $UAR testing-migrations
# Debian bugs
0 */4 * * * $UAR bugs 2>&1 | /org/udd.debian.org/udd/scripts/filter-bugs-output.sh
0 13 * * * $UAR bugs-archive 2>&1 | /org/udd.debian.org/udd/scripts/filter-bugs-output.sh

# Ubuntu bugs
0 4 * * * $UAR ubuntu-bugs

# Dump DB
0 8 */2 * * /org/udd.debian.org/udd/scripts/dump-db.sh
