# crontabs
# m h  dom mon dow   command
MAILTO=lucas@debian.org
UAR=/org/udd.debian.org/udd/update-and-run.sh
# Debian Sources/Packages
0 */12 * * * $UAR debian-squeeze debian-sid debian-experimental
0 3 * * * $UAR debian-etch debian-backports-etch debian-volatile-etch debian-etch-proposed-updates debian-etch-security
0 15 * * * $UAR debian-lenny debian-backports-lenny debian-volatile-lenny debian-lenny-proposed-updates debian-lenny-security debian-squeeze-proposed-updates debian-squeeze-security
# Ubuntu Sources/Packages
30 2 * * * $UAR ubuntu-lucid ubuntu-karmic ubuntu-hardy ubuntu-intrepid ubuntu-jaunty
# Various simple things
0 4,16 * * * $UAR dehs debian-popcon ubuntu-popcon lintian debtags carnivore ldap removals | /org/udd.debian.org/udd/scripts/filter-output.rb
49 */6 * * * $UAR upload-history | /org/udd.debian.org/udd/scripts/filter-output.rb
0 */12 * * * $UAR testing-migrations
11 */4 * * * $UAR wannabuild | /org/udd.debian.org/udd/scripts/filter-output.rb
# Debian bugs
0 1,7,13,19 * * * $UAR bugs orphaned-packages 2>&1 | /org/udd.debian.org/udd/scripts/filter-bugs-output.sh | /org/udd.debian.org/udd/scripts/filter-output.rb
30 14 * * * $UAR bugs-archive 2>&1 | /org/udd.debian.org/udd/scripts/filter-bugs-output.sh

# Ubuntu bugs
57 9 * * * $UAR ubuntu-bugs

# Dump DB
0 8 */2 * * /org/udd.debian.org/udd/scripts/dump-db.sh
0 4 */2 * * /org/udd.debian.org/udd/scripts/generate-schema.sh

# checks
*/5 * * * * /org/udd.debian.org/udd/scripts/monitor_locks
0 * * * * /org/udd.debian.org/udd/scripts/check_timestamps
0 * * * * crontab -l > /org/udd.debian.org/udd/web/crontabs.txt
# test script for the removal gatherer
0 1 2 * * python /org/udd.debian.org/udd/udd/removals_gatherer.py /org/udd.debian.org/mirrors/removals-full.txt > /dev/null

MAILTO=tille@debian.org
# screenshots.debian.net
30 3 * * * $UAR screenshots | /org/udd.debian.org/udd/scripts/filter-output.rb 2>&1 | mail -e -s 'UDD screenshots cron job' tille@debian.org

# DDTP
0 2 * * * $UAR ddtp | /org/udd.debian.org/udd/scripts/filter-output.rb 2>&1 | mail -e -s 'UDD ftpnew cron job' tille@debian.org

# FTP new queue
0 3 * * * $UAR ftpnew | /org/udd.debian.org/udd/scripts/filter-output.rb 2>&1 | mail -e -s 'UDD ftpnew cron job' tille@debian.org
